# enterprise/bw-enterprise-backup.timer
[Unit]
Description=Daily BW Monitor v49 enterprise backup

[Timer]
OnCalendar=*-*-* 04:10:00
RandomizedDelaySec=20m
Persistent=true
Unit=bw-enterprise-backup.service

[Install]
WantedBy=timers.target
